# deploy/farmore-resolver.service
[Unit]
Description=Farmore resolver (handle resolution / send-to-handle / sign-in) — read-only HTTP service
Documentation=https://github.com/farmore-network/farmore-node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=farmore
Group=farmore
EnvironmentFile=/etc/farmore/resolver.env
ExecStart=/opt/farmore/bin/farmore-resolver
Restart=on-failure
RestartSec=5
# Start-rate limiting so a crash loop backs off instead of hammering.
StartLimitIntervalSec=60
StartLimitBurst=5

# --- Hardening (read-only service, holds no funds) ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
